# lec-7233/etc/uci-defaults/99-custom
# settings

uci set.network.lan.proto="dhcp"
uci commit network

# timezone

uci set system.@system[0].hostname='OpenWrt'
uci set system.@system[0].zonename='Europe/Moscow'
uci set system.@system[0].timezone='MSK-3'
uci commit system

# ntp servers

uci set system.@timeserver[0].enabled='1'
uci set system.@timeserver[0].enable_server='0'
uci set system.@timeserver[0].use_dhcp='0'
uci delete system.@timeserver[0].server
uci add_list system.@timeserver[0].server='0.ru.pool.ntp.org'
uci add_list system.@timeserver[0].server='1.ru.pool.ntp.org'
uci add_list system.@timeserver[0].server='2.ru.pool.ntp.org'
uci add_list system.@timeserver[0].server='3.ru.pool.ntp.org'
uci commit system
